# deploy/balam-opencode.service
[Unit]
Description=OpenCode server (Balam agent backend, ADR-0001)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
# OPENCODE_SERVER_PASSWORD (and OPENCODE_SERVER_USERNAME) come from the repo .env;
# opencode serve reads them from its environment. Bound to localhost (ADR-0007).
EnvironmentFile=/home/ubuntu/projects/balam/.env
# Mirror the interactive zsh PATH so OpenCode's non-interactive tool shells can
# find user-installed CLIs (glab, lark-cli, acli, iris, tolgee, figma, node,
# bun, Go-installed binaries). systemd does not source .zshrc/.zshenv.
# The nvm version is pinned because nvm.sh cannot run in a unit; bump it when
# switching the default Node version.
Environment=PATH=/home/ubuntu/.bun/bin:/home/ubuntu/.opencode/bin:/home/ubuntu/.local/bin:/home/ubuntu/.nvm/versions/node/v24.15.0/bin:/home/ubuntu/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/usr/local/go/bin:/home/ubuntu/go/bin
Environment=GOPATH=/home/ubuntu/go
# Plan mode is experimental and gated off by default on the server (registry.ts:
# `flags.experimentalPlanMode && flags.client === "cli"`; client defaults to
# "cli"). This exposes plan_enter/plan_exit to the model so Telegram sessions
# can enter plan mode; permissions.py allows both in the session ruleset.
Environment=OPENCODE_EXPERIMENTAL_PLAN_MODE=true
# Dedicated SQLite database. opencode migrates its DB forward on every boot
# and migrations are one-way, so on the default shared DB any newer binary
# booting anywhere on this VM (ivy, interactive runs) breaks this
# still-running server mid-flight (2026-06-12: "NOT NULL constraint failed:
# session_message.seq" on every prompt). Auth/logs/diffs stay shared.
Environment=OPENCODE_DB=/home/ubuntu/.local/share/opencode/balam.db
ExecStart=/home/ubuntu/.opencode/bin/opencode serve --hostname 127.0.0.1 --port 4096
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
